# rekishiya-config/dynamic-dns.timer
[Unit]
Description=Periodically trigger dynamic DNS update
After=network.target

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
